# env-v2/env.service
[Unit]
Description=Temperature monitoring display for multiple Raspberry Pi computers.

[Service]
User=pi
Type=idle
ExecStart=/home/pi/env.py

[Install]
WantedBy=multi-user.target